# systemd/simpl-local-usersroles.service
[Unit]
Description=Simpl onboarding env=local

[Service]
ExecStart=mise --cd "${SIMPL_PROJECT_HOME}" run run:users-roles
StandardOutput=journal
StandardError=journal
